# packaging/systemd/system/cerebrix-ingest-journald.service
[Unit]
Description=Halbert Ingestion (journald)
After=network-online.target systemd-journald.service
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
SupplementaryGroups=systemd-journal
StateDirectory=halbert
LogsDirectory=halbert
ConfigurationDirectory=halbert
Environment=Halbert_DATA_DIR=/var/lib/halbert
Environment=Halbert_LOG_DIR=/var/log/halbert
# Adjust ExecStart to your install path or packaged CLI
ExecStart=/usr/bin/env python3 /opt/halbert/Halbert/main.py ingest-journald
Restart=on-failure
RestartSec=2s
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX
SystemCallFilter=@system-service
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
