# packaging/systemd/system/cerebrix-ingest-journald.service
[Unit]
Description=Cerebric Ingestion (journald)
After=network-online.target systemd-journald.service
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
SupplementaryGroups=systemd-journal
StateDirectory=cerebric
LogsDirectory=cerebric
ConfigurationDirectory=cerebric
Environment=Cerebric_DATA_DIR=/var/lib/cerebric
Environment=Cerebric_LOG_DIR=/var/log/cerebric
# Adjust ExecStart to your install path or packaged CLI
ExecStart=/usr/bin/env python3 /opt/cerebric/Cerebric/main.py ingest-journald
Restart=on-failure
RestartSec=2s
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX
SystemCallFilter=@system-service
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
