# packaging/systemd/system/halbert-dashboard.service
[Unit]
Description=Halbert Dashboard API Server
Documentation=https://github.com/halbert/LinuxBrain
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/opt/halbert
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/%i

# Use the virtual environment
ExecStart=/opt/halbert/.venv/bin/python -m uvicorn \
    halbert_core.halbert_core.dashboard.app:app \
    --host 127.0.0.1 \
    --port 8000

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/opt/halbert/data /tmp

[Install]
WantedBy=multi-user.target
