# packaging/systemd/system/halbert-config-watch.path
[Unit]
Description=Path trigger for Halbert Config Watcher

[Path]
# Trigger when files under /etc/halbert change (system scope)
PathModified=/etc/halbert/
Unit=halbert-config-watch.service

[Install]
WantedBy=multi-user.target
